# crontab
# m h  dom mon dow   command
* * * * * /root/ElCheapoAIS-manhole/manhole.sh /etc/elcheapoais
