# crontab
# m h dom mon dow user  command
  * * *   *   *   root  /usr/local/bin/manhole.sh /etc/elcheapoais > /tmp/manhole.txt 2>&1
